# flaskblog.conf
[program:flaskblog]
command = gunicorn -b 0.0.0.0:8005 wsgi_gunicorn:app
directory = /home/xin/www/flaskblog
user = xin