# flaskblog.conf
[program:flaskblog]
command = gunicorn -b 0.0.0.0:8005 run:app
directory = /home/xin/www/flaskblog
user = xin