# flaskblog.conf
[program:flaskblog]
command = gunicorn -b 127.0.0.1:9000 wsgi_gunicorn:app
directory = /home/xin/workspace/python/flaskblog
user = xin
