# etc/flaskblog.conf
[program:flaskblog]
command = /home/xin/workspace/python/flaskblog/venv/bin/gunicorn -b 127.0.0.1:9000 wsgi_gunicorn:app
directory = /home/xin/workspace/python/flaskblog
user = xin
